# infra/systemd/battleship-backup.service
[Unit]
Description=BattleShipArena SQLite snapshot backup
After=network-online.target

[Service]
Type=oneshot
User=battleship
Group=battleship
WorkingDirectory=/opt/battleship-arena
ExecStart=/usr/bin/env bash /opt/battleship-arena/infra/scripts/backup.sh
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
MemoryDenyWriteExecute=yes
ReadWritePaths=/var/lib/battleship-arena /var/lib/battleship-arena-staging /var/backups/battleship-arena
